# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:inference]
command=python3 -m uvicorn ml.serve.inference_app:app --host %(ENV_INFERENCE_BIND_HOST)s --port 8000
directory=/app
autorestart=true
startretries=10
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=CROPINTEL_BACKEND="tflite",CROPINTEL_PREDICTION_LOG="/app/data/predictions.jsonl"

[program:web]
command=node node_modules/.bin/next start -H 0.0.0.0 -p 3050
directory=/app
autorestart=true
startretries=10
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",INFERENCE_URL="http://127.0.0.1:8000"
